# deploy/systemd/ipright-api.service
[Unit]
Description=IPRight FastAPI service
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/opt/ipright/backend
EnvironmentFile=/opt/ipright/backend/.env.production
ExecStart=/opt/ipright/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 18000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
