# deploy/systemd/ipright-worker.service
[Unit]
Description=IPRight Celery worker
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/opt/ipright/current
EnvironmentFile=/opt/ipright/backend/.env.production
Environment=PYTHONPATH=/opt/ipright/current/backend:/opt/ipright/current
Environment=PLAYWRIGHT_BROWSERS_PATH=/opt/ipright/shared/ms-playwright
Environment=TMPDIR=/opt/ipright/shared/tmp
ExecStart=/opt/ipright/backend/.venv/bin/celery -A workers.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
